# CFTR/mailsender.service
[Unit]
Description=mailsender

[Service]
Type=simple
StandardOutput=journal
ExecStart=/opt/nys_nbs/mailsender.sh

[Install]
WantedBy=default.target